# deploy/mitayshvim.service
[Unit]
Description=Small Mitayshvim game server
After=network.target

[Service]
User=mitayshvim
WorkingDirectory=/opt/mitayshvim
ExecStart=/opt/mitayshvim/mitayshvim -addr 127.0.0.1:8080 -data /opt/mitayshvim/data
Restart=on-failure
TimeoutStopSec=15

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mitayshvim/data
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
